# deploy/systemd/findata-backup.timer
[Unit]
Description=FindataChile: respaldo diario a Google Drive

[Timer]
OnCalendar=*-*-* 07:30
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
